# launchagents/Library/LaunchAgents/com.brew.update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
	<dict>
		<key>Label</key>
		<string>com.brew.update</string>
		<key>ProgramArguments</key>
		<array>
			<string>/bin/bash</string>
			<string>-c</string>
			<string>
				/opt/homebrew/bin/brew update
				/opt/homebrew/bin/brew upgrade
				/opt/homebrew/bin/brew bundle check || /opt/homebrew/bin/brew bundle install
				/opt/homebrew/bin/brew services info logrotate || /opt/homebrew/bin/brew services restart logrotate
			</string>
		</array>
		<key>KeepAlive</key>
		<dict>
			<key>SuccessfulExit</key>
			<false/>
		</dict>
		<key>StandardErrorPath</key>
		<string>/opt/homebrew/var/log/brewupdateError.log</string>
		<key>StandardOutputPath</key>
		<string>/opt/homebrew/var/log/brewupdate.log</string>
	</dict>
</plist>
