# launchagents/Library/LaunchAgents/com.brew.upgrade.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
	<dict>
		<key>Label</key>
		<string>com.brew.upgrade</string>
		<key>ProgramArguments</key>
		<array>
			<string>/bin/bash</string>
			<string>-c</string>
			<string>/opt/homebrew/bin/brew</string>
			<string>upgrade</string>
		</array>
		<key>RunAtLoad</key>
		<true/>
		<key>StartInterval</key>
		<integer>7200</integer>
		<key>StandardErrorPath</key>
		<string>/opt/homebrew/var/log/upgrade.log</string>
		<key>StandardOutPath</key>e
		<string>/opt/homebrew/var/log/upgrade.log</string>
	</dict>
</plist>
