# launchagents/Library/LaunchAgents/com.dotfiles.update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
	<dict>
		<key>Label</key>
		<string>com.dotfiles.update</string>
		<key>ProgramArguments</key>
		<array>
			<string>/bin/bash</string>
			<string>-c</string>
			<string>
				if [ -d $HOME/github/dotfiles ]; then
					pushd $HOME/github/dotfiles
					git pull --recurse-submodules
					popd
				fi
			</string>
		</array>
		<key>KeepAlive</key>
		<dict>
			<key>SuccessfulExit</key>
			<false/>
		</dict>
		<key>StandardErrorPath</key>
		<string>/opt/homebrew/var/log/dotfilesupdate.log</string>
		<key>StandardOutputPath</key>
		<string>/opt/homebrew/var/log/dotfilesupdate.log</string>
	</dict>
</plist>
